# daemon/co.welf.amber-cool.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.welf.amber-cool</string>

    <!-- Runs the fan-control daemon as root; it reads its mode from the config file. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/fanctl</string>
        <string>daemon</string>
        <string>/usr/local/etc/amber-cool/mode</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it dies; full blast (or the chosen mode) stays engaged. -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/amber-cool.log</string>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/amber-cool.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
